# deploy/selfhosted/trading-bot.service
# Systemd Service para Trading Bot
# Instalar em: /etc/systemd/system/trading-bot.service
#
# Comandos:
#   sudo systemctl daemon-reload
#   sudo systemctl enable trading-bot.timer
#   sudo systemctl start trading-bot.timer
#   sudo systemctl status trading-bot
#   journalctl -u trading-bot -f

[Unit]
Description=Intelligent Trading Bot
After=network.target

[Service]
Type=oneshot
User=pi
Group=pi

# Diretório do bot
WorkingDirectory=/home/pi/trading-bot

# Variáveis de ambiente
Environment="TRADING_MODE=paper"
Environment="BOT_DATA_DIR=/home/pi/trading-bot/data"

# Se usar Telegram
# Environment="TELEGRAM_BOT_TOKEN=seu_token"
# Environment="TELEGRAM_CHAT_ID=seu_chat_id"

# Se usar exchange (modo live)
# Environment="EXCHANGE_API_KEY=sua_key"
# Environment="EXCHANGE_SECRET=seu_secret"

# Executar
ExecStart=/usr/bin/python3 /home/pi/trading-bot/deploy/selfhosted/run_daemon.py

# Logs
StandardOutput=append:/home/pi/trading-bot/logs/systemd.log
StandardError=append:/home/pi/trading-bot/logs/systemd_error.log

# Restart em caso de falha
Restart=on-failure
RestartSec=60

[Install]
WantedBy=multi-user.target
